# mcp-server/wireguard/wg-reload.path
[Unit]
Description=Watch for WireGuard peer config changes

[Path]
PathModified=/opt/mcp-server/wireguard/peers

[Install]
WantedBy=multi-user.target
